# .config/autostart/pipewire-pulse.desktop
[Desktop Entry]
Name=pipewire-pulse
Comment=Pipewire-pulse bridge
Exec=pipewire-pulse
Type=Application
